# prevent-dpms.service
[Unit]
Description=Prevent DPMS when sound is playing

[Service]
Type=simple
ExecStart=/usr/bin/prevent-dpms

[Install]
WantedBy=default.target
